# coffee-duty.service
[Unit]
Description=Coffee Duty Laravel Application
After=network.target

[Service]
Type=simple
User=cakicitu
Group=cakicitu
WorkingDirectory=/opt/projects/coffee-duty
ExecStart=/usr/bin/php artisan serve --host=0.0.0.0 --port=8001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
